# graccbackup@.service
[Unit]
Description=GRACC Periodic Backup
Documentation=https://opensciencegrid.github.io/gracc
Wants=network-online.target
After=network-online.target

[Service]
User=gracc
Group=gracc
Type=oneshot
Environment=IN=%i
PermissionsStartOnly=true
ExecStart=sudo /bin/docker run -v /etc:/etc -v /data/graccarchive:/var/lib/graccarchive gracc-backup %i
TimeoutStopSec=25min

[Install]
WantedBy=multi-user.target
