# gracc-backup-ps-itb.service
[Unit]
Description=GRACC Periodic Backup for PS ITB
Documentation=https://opensciencegrid.github.io/gracc
Wants=network-online.target
After=network-online.target

[Service]
User=gracc
Group=gracc
Type=oneshot
ExecStart=/usr/bin/gracc-backup /var/lib/psarchive-itb/output /var/lib/psarchive-itb/secondary gsiftp://fndca1.fnal.gov/pnfs/fs/usr/fermigrid/gratia/gracc-ps-itb-raw/
TimeoutStopSec=5min

[Install]
WantedBy=multi-user.target

